# svc/controller-system/make-deb/debian/postrm
#!/bin/sh

( systemctl stop eva-cs-agent > /dev/null 2>&1 || exit 0 )
( systemctl disable eva-cs-agent > /dev/null 2>&1 || exit 0 )
( userdel eva-cs-agent > /dev/null 2>&1 || exit 0 )
